# tpl_whisper.socket
[Unit]
Description=Whisper Socket
PartOf=whisper.service

[Socket]
ListenStream=${WASTE_SOCKET}

[Install]
WantedBy=sockets.target
