# tpl_waste.socket
[Unit]
Description=Waste Socket
PartOf=waste.service

[Socket]
ListenStream=${WASTE_SOCKET}

[Install]
WantedBy=sockets.target
